# systemd/wikisync.service
[Unit]
Description=Wikipedia Local Synchronization Service
Documentation=https://github.com/martindale/wikisync
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/wikipedia
ExecStart=/usr/bin/python3 /opt/wikipedia/wikisync.py --service
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wikisync

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/wikipedia

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
